# pac4cli.service
[Unit]
Description=PAC autoconfigured proxy for use through http_proxy= environment variables
After=network.target

[Service]
ExecStart=/usr/local/bin/pac4cli -p 3128 --systemd --config /etc/pac4cli/pac4cli.config --loglevel warn
Restart=always

[Install]
WantedBy=network.target
