# pac4cli.service
[Unit]
Description=PAC autoconfigured proxy for use through http_proxy= environment variables
After=network.target

[Service]
ExecStart=/usr/bin/pac4cli -p 3128 --systemd --loglevel warn
Restart=always

[Install]
WantedBy=network.target
