# debian/postinst
#!/bin/sh
# postinst script for pac4cli
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

start_service() {
    echo "Starting pac4cli service..."
    systemctl enable pac4cli.service
    systemctl start pac4cli.service
}


case "$1" in
    configure)
        # this old configuration file breaks pac4cli startup
        # because it sends SIGHUPs that we don't handle anymore
        # since e2f264aa5adf7f4a3c2d3ed2. Technically, we shouldn't
        # just remove this but check for local modifications, but
        # that's sufficiently unlikelye that I'll just go with this.
        rm -f /etc/NetworkManager/dispatcher.d/trigger-pac4cli
        start_service
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
