# debian/prerm
#!/bin/sh
echo "Stopping pac4cli.service..."
systemctl stop pac4cli.service
systemctl disable pac4cli.service
